# docker/backup.crontab
*/1 * * * * /usr/bin/sh /var/configs/backup.sh
