# docker/backup.crontab
0 3 * * * /usr/bin/sh /var/configs/backup.sh
